# smart-greeter.service
[Unit]
Description=Smart AI Boot Greeter
After=network.target sound.target graphical-session.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/shelemiah/Personal/os_ai
Environment=PATH=/home/shelemiah/Personal/os_ai/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/shelemiah/Personal/os_ai/venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/shelemiah/Personal/os_ai
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/shelemiah/Personal/os_ai

[Install]
WantedBy=default.target
